# base-deb10/linux-config/etc/direct/profile.d/oinkserver.sh
#!/bin/bash

# Local Consul CLI config
export CONSUL_HTTP_ADDR="169.254.1.1:8500"
export CONSUL_RPC_ADDR="169.254.1.1:8400"
